# com.patterniq.daily.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.patterniq.daily</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/tamirreznik/code/private/PatternIQ/venv/bin/python</string>
        <string>/Users/tamirreznik/code/private/PatternIQ/macos_daily_runner.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/tamirreznik/code/private/PatternIQ</string>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>18</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/tamirreznik/code/private/PatternIQ/logs/patterniq_daily.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/tamirreznik/code/private/PatternIQ/logs/patterniq_daily_error.log</string>

    <key>RunAtLoad</key>
    <false/>

    <key>KeepAlive</key>
    <false/>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/tamirreznik/code/private/PatternIQ/venv/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>PYTHONPATH</key>
        <string>/Users/tamirreznik/code/private/PatternIQ</string>
    </dict>
</dict>
</plist>
